# docker/config/supervisor/conf.d/deliversm.ini
[program:deliversmd]
command=python -u /usr/local/bin/deliversmd.py
autostart=true
autorestart=true