# docker/config/supervisor/conf.d/dlrlookup.ini
[program:dlrlookupd]
command=python -u /usr/local/bin/dlrlookupd.py
autostart=true
autorestart=true