# docker/config/supervisor/conf.d/dlrthrower.ini
[program:dlrd]
command=python -u /usr/local/bin/dlrd.py
autostart=true
autorestart=true