# docker/config/supervisor/conf.d/interceptor.ini
[program:interceptord]
command=python -u /usr/local/bin/interceptord.py
autostart=true
autorestart=true