# docker/config/supervisor/conf.d/jasmin.ini
[program:jasmind]
command=python -u /usr/local/bin/jasmind.py --enable-interceptor-client -u jcliadmin -p jclipwd
autostart=true
autorestart=true
