# docker/config/supervisor/conf.d/monitor.ini
[program:monitor]
command=/etc/jasmin/scripts/monitor.sh
autostart=true
autorestart=true
